# files/zfs/usr/lib/systemd/system/ensure-sysext.service.d/10-load-modules.conf
[Service]
ExecStartPost=/usr/bin/systemctl restart systemd-modules-load.service
