# nvidia-runtime.sysext/files/usr/lib/systemd/system/nvidia.service.d/10-persistenced.conf
[Service]
ExecStartPre=/usr/local/bin/_nvidia-runtime-helper pre
ExecStartPost=/usr/local/bin/_nvidia-runtime-helper post
ExecStopPost=rmmod nvidia_uvm nvidia_modeset nvidia
